# archive/infra/ufobeep-api.service

[Unit]
Description=UFOBeep FastAPI Service
After=network.target

[Service]
User=ufobeep
Group=ufobeep
WorkingDirectory=/home/ufobeep/ufobeep/api
Environment="PATH=/home/ufobeep/ufobeep/venv/bin"
ExecStart=/home/ufobeep/ufobeep/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
